# coordinador/com.sergi.aura-coordinador.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sergi.aura-coordinador</string>

  <!-- S'obre via Terminal (open) perquè launchd no pot llegir ~/Documents -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/open</string>
    <string>-g</string>
    <string>/Users/sergicastillo/Documents/Aura/coordinador/coordinador.command</string>
  </array>

  <!-- Cada dia a les 9:00, després del briefing (07:32) i l'Escriptor (08:06) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
